# deploy/utility-jkh.service.example
# Шаблон unit-файла. Подставьте свой путь и пользователя.
# Установка:
#   sudo cp deploy/utility-jkh.service /etc/systemd/system/
#   sudo sed -i 's|/opt/utility-jkh|/ваш/путь|g' /etc/systemd/system/utility-jkh.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now utility-jkh
#   sudo systemctl status utility-jkh
#   journalctl -u utility-jkh -f

[Unit]
Description=Система учета ЖКХ (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/utility-jkh
EnvironmentFile=-/opt/utility-jkh/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/utility-jkh
ExecStart=/opt/utility-jkh/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=utility-jkh

[Install]
WantedBy=multi-user.target
